# provisioning/files/etc/update-motd.d/manati-terminal-text
[48;5;22m  [48;5;28m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;28m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;23m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;22m  [48;5;16m  [48;5;22m  [48;5;28m  [48;5;28m  [48;5;28m  [48;5;28m  [48;5;28m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;22m  [0m
[48;5;71m  [48;5;22m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;22m  [48;5;71m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;35m  [48;5;22m  [48;5;29m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;34m  [48;5;65m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;35m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;65m  [48;5;22m  [48;5;77m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;65m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;29m  [48;5;22m  [0m
[48;5;71m  [48;5;16m  [48;5;23m  [48;5;16m  [48;5;16m  [48;5;29m  [48;5;22m  [48;5;71m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;22m  [48;5;16m  [48;5;28m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;34m  [48;5;16m  [48;5;70m  [48;5;16m  [48;5;16m  [48;5;34m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;65m  [48;5;16m  [48;5;22m  [48;5;28m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;16m  [0m
[48;5;71m  [48;5;16m  [48;5;29m  [48;5;16m  [48;5;16m  [48;5;35m  [48;5;16m  [48;5;71m  [48;5;16m  [48;5;16m  [48;5;34m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;35m  [48;5;16m  [48;5;71m  [48;5;22m  [48;5;16m  [48;5;34m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;34m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;28m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;28m  [48;5;16m  [0m
[48;5;71m  [48;5;16m  [48;5;29m  [48;5;16m  [48;5;22m  [48;5;22m  [48;5;16m  [48;5;71m  [48;5;16m  [48;5;16m  [48;5;65m  [48;5;16m  [48;5;16m  [48;5;23m  [48;5;28m  [48;5;16m  [48;5;16m  [48;5;34m  [48;5;16m  [48;5;16m  [48;5;77m  [48;5;16m  [48;5;34m  [48;5;16m  [48;5;16m  [48;5;28m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;65m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;28m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;28m  [48;5;16m  [0m
[48;5;71m  [48;5;16m  [48;5;16m  [48;5;71m  [48;5;71m  [48;5;16m  [48;5;16m  [48;5;71m  [48;5;16m  [48;5;22m  [48;5;34m  [48;5;28m  [48;5;34m  [48;5;28m  [48;5;35m  [48;5;16m  [48;5;16m  [48;5;34m  [48;5;16m  [48;5;16m  [48;5;28m  [48;5;22m  [48;5;28m  [48;5;22m  [48;5;16m  [48;5;28m  [48;5;28m  [48;5;34m  [48;5;28m  [48;5;35m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;28m  [48;5;16m  [0m
[48;5;71m  [48;5;16m  [48;5;16m  [48;5;41m  [48;5;71m  [48;5;16m  [48;5;16m  [48;5;71m  [48;5;16m  [48;5;29m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;28m  [48;5;22m  [48;5;16m  [48;5;34m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;64m  [48;5;77m  [48;5;16m  [48;5;22m  [48;5;29m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;29m  [48;5;16m  [48;5;16m  [48;5;28m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;28m  [48;5;16m  [0m
[48;5;23m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;23m  [48;5;16m  [48;5;23m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;23m  [48;5;16m  [48;5;59m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;29m  [48;5;16m  [48;5;22m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;23m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;22m  [48;5;16m  [48;5;16m  [48;5;16m  [48;5;22m  [48;5;16m  [0m

